# _deployment/systemd/intermediate-app.service
[Unit]
Description=Intermediate Service Laravel App
After=network.target redis.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/html/AIScreen-AICamera/intermediate-service
Environment="PATH=/usr/bin:/usr/sbin"
ExecStart=/usr/bin/php artisan serve --host=0.0.0.0 --port=8080
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

